# hc/alerts/root/etc/services.d/alerts/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid nobody "${APP_PATH}"/manage.py sendalerts
