# app/root/etc/services.d/prunechecks/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid nobody snooze -v -d1 "${APP_PATH}"/manage.py prunenotifications 2>&1 | sed --unbuffered 's/.*/prunenotifications: &/'
